# deploy/ewc-live.service
[Unit]
Description=EWC Live App
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ewc-live-app
EnvironmentFile=/opt/ewc-live-app/.env
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=2

# SSE needs many open sockets
LimitNOFILE=65535

# Helpful defaults
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target

